# systemd/pawns-watchdog.timer
[Unit]
Description=Run pawns-cli watchdog every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
